# Hunting/Hunt-RemcosRAT-JScript-JPEG-Dropper.json
{
  "id": "03f4ad08-0539-44dc-b801-3378d1690a58",
  "name": "Hunt: Remcos RAT JScript/JPEG Multi-Stage Dropper Campaign",
  "description": "Hunting query for the multi-stage JScript → PowerShell → JPEG payload → Remcos RAT campaign documented by SANS ISC (2026-02-23). The attack uses heavily obfuscated JScript attached to phishing emails (spoofed sender, embedded logo). The JScript achieves startup folder persistence, spawns a hidden PowerShell via WMI Win32_Process.Create to decode Base64 payload, downloads a JPEG image from archive.org with embedded PE payload, and ultimately delivers Remcos RAT from a secondary C2. Error messages appear in Brazilian Portuguese.\n\nIOC SHA256 Hashes:\n- JS 1st stage: edc04c2ab377741ef50b5ecbfc90645870ed753db8a43aa4d0ddcd26205ca2a4\n- Remcos RAT EXE: 1158ef7830d20d6b811df3f6e4d21d41c4242455e964bde888cd5d891e2844da\n- 3rd stage encoded TXT: b6fdb00270914cdbc248cacfac85749fa7445fca1122a854dce7dea8f251019c\n\nKnown C2 URLs (defanged):\n- hxxps[:]//hotelseneca[.]ro/ConvertedFileNew.txt\n- hxxps[:]//ia600603[.]us[.]archive[.]org/13/items/msi-pro-with-b-64_202602/MSI_PRO_with_b64.png\n\nReference: https://isc.sans.edu/diary.html (2026-02-23)\nRelated diaries: https://isc.sans.edu/diary/32682 https://isc.sans.edu/diary/32726",
  "category": "Hunting Queries",
  "version": 2,
  "tags": [
    { "name": "tactics", "value": "InitialAccess,Execution,Persistence,DefenseEvasion,CommandAndControl" },
    { "name": "techniques", "value": "T1566,T1059,T1547,T1027,T1071" },
    { "name": "source", "value": "SANS ISC 2026-02-23" },
    { "name": "malware", "value": "RemcosRAT" },
    { "name": "severity", "value": "High" }
  ],
  "query": "// Hunt: Remcos RAT JScript/JPEG Multi-Stage Dropper (SANS ISC 2026-02-23)\n// Phase 1: Known IOC file hashes\nlet KnownHashes = dynamic([\n    'edc04c2ab377741ef50b5ecbfc90645870ed753db8a43aa4d0ddcd26205ca2a4',  // JS 1st stage\n    '1158ef7830d20d6b811df3f6e4d21d41c4242455e964bde888cd5d891e2844da',  // Remcos RAT EXE\n    'b6fdb00270914cdbc248cacfac85749fa7445fca1122a854dce7dea8f251019c'   // 3rd stage encoded TXT\n]);\n// Phase 2: Known C2 domains\nlet KnownC2Domains = dynamic(['hotelseneca.ro', 'archive.org']);\nlet KnownC2Paths = dynamic(['/ConvertedFileNew.txt', 'MSI_PRO_with_b64.png', 'msi-pro-with-b-64']);\n// IOC Match: file hashes\nlet HashMatches = union isfuzzy=true (\n    DeviceFileEvents\n    | where TimeGenerated > ago(30d)\n    | where SHA256 in (KnownHashes)\n    | extend MatchType = \"Known IOC Hash\", DeviceName = DeviceName, Detail = FileName\n    | project TimeGenerated, DeviceName, MatchType, Detail, SHA256, InitiatingProcessAccountName\n),\n(\n    DeviceProcessEvents\n    | where TimeGenerated > ago(30d)\n    | where SHA256 in (KnownHashes)\n    | extend MatchType = \"Known IOC Hash (Process)\", Detail = FileName, InitiatingProcessAccountName = InitiatingProcessAccountName\n    | project TimeGenerated, DeviceName, MatchType, Detail, SHA256, InitiatingProcessAccountName\n);\n// Behavioral: WMI spawning hidden PowerShell with Invoke-Expression / Base64 decode\nlet WMIPowerShell = DeviceProcessEvents\n    | where TimeGenerated > ago(30d)\n    | where InitiatingProcessFileName =~ \"WmiPrvSE.exe\"\n      and FileName =~ \"powershell.exe\"\n      and (ProcessCommandLine has_any (\"Invoke-Expression\", \"IEX\", \"EncodedCommand\", \"-enc \", \"FromBase64String\", \"UTF-8\"))\n    | extend MatchType = \"WMI-spawned PowerShell with Base64 decode (JScript dropper TTP)\", Detail = ProcessCommandLine, SHA256 = \"\"\n    | project TimeGenerated, DeviceName, MatchType, Detail, SHA256, InitiatingProcessAccountName;\n// Behavioral: JScript file executing from non-standard path with large file size indicator\nlet JScriptExec = DeviceProcessEvents\n    | where TimeGenerated > ago(30d)\n    | where FileName in~ (\"wscript.exe\", \"cscript.exe\")\n      and ProcessCommandLine has_any (\".js\", \".jse\")\n      and ProcessCommandLine !has_any (\"C:\\\\Program Files\\\\\", \"C:\\\\Windows\\\\\", \"C:\\\\Program Files (x86)\\\\\")\n    | extend MatchType = \"JScript executed from non-standard path\", Detail = ProcessCommandLine, SHA256 = \"\"\n    | project TimeGenerated, DeviceName, MatchType, Detail, SHA256, InitiatingProcessAccountName;\n// Behavioral: Startup folder persistence via file copy from script\nlet StartupPersistence = DeviceFileEvents\n    | where TimeGenerated > ago(30d)\n    | where FolderPath has_any (\"\\\\AppData\\\\Roaming\\\\Microsoft\\\\Windows\\\\Start Menu\\\\Programs\\\\Startup\\\\\", \"\\\\ProgramData\\\\Microsoft\\\\Windows\\\\Start Menu\\\\Programs\\\\StartUp\\\\\")\n      and FileName has_any (\".js\", \".jse\", \".vbs\", \".vbe\", \".wsf\")\n      and ActionType in~ (\"FileCreated\", \"FileModified\")\n    | extend MatchType = \"Script dropped in Startup folder (JScript persistence TTP)\", Detail = FolderPath, SHA256 = \"\"\n    | project TimeGenerated, DeviceName, MatchType, Detail, SHA256, InitiatingProcessAccountName;\n// Network: C2 domains\nlet C2Network = DeviceNetworkEvents\n    | where TimeGenerated > ago(30d)\n    | where RemoteUrl has_any (KnownC2Paths)\n       or (RemoteUrl has \"hotelseneca.ro\" and RemoteUrl has \".txt\")\n    | extend MatchType = \"Connection to known C2 URL\", Detail = RemoteUrl, SHA256 = \"\"\n    | project TimeGenerated, DeviceName, MatchType, Detail, SHA256, InitiatingProcessAccountName = InitiatingProcessAccountName;\n// Union all\nunion isfuzzy=true HashMatches, WMIPowerShell, JScriptExec, StartupPersistence, C2Network\n| summarize Count = count(), EarliestSeen = min(TimeGenerated), LatestSeen = max(TimeGenerated), Details = make_set(Detail, 5) by DeviceName, MatchType, SHA256\n| sort by LatestSeen desc"
}
